# user/unclutter.service
[Unit]
Description = Unclutter

[Service]
ExecStart = /usr/bin/unclutter -idle 3 -root

[Install]
WantedBy = default.target
